# deploy/drs-frontend.service
# /etc/systemd/system/drs-frontend.service
# sudo systemctl daemon-reload && sudo systemctl enable --now drs-frontend

[Unit]
Description=DRS Frontend (static build server)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=alis
WorkingDirectory=/home/alis/repos/drs-llm/frontend
Environment="NODE_ENV=production" "PORT=3173"
ExecStart=/home/alis/.nvm/versions/node/v22.19.0/bin/npm run start
Restart=always
RestartSec=5
# optional hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
